# contrib/configurations/development/supervisord.conf
[program:project_core]
command = ./project_core.py start
directory = /home/project_core/
user = project_core
autostart = true
autorestart = true
redirect_stderr = True
priority = 10
